# configs/etc/systemd/system/crash-reporter.service
[Unit]
Description=System Crash Reporter
After=network.target kernel-logger.service
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/crash_reporter.py
RemainAfterExit=yes

# Security settings
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadOnlyPaths=/proc/sys/kernel/random/boot_id /proc/uptime /proc/meminfo /proc/modules
CapabilityBoundingSet=CAP_DAC_READ_SEARCH

[Install]
WantedBy=multi-user.target