# configs/etc/systemd/system/kernel-logger.service
[Unit]
Description=Kernel Logger Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/kernel_logger_daemon.py
Restart=always
RestartSec=5
Environment=KERNEL_LOG_PATH=/var/log/kernel
Environment=KERNEL_LOG_LEVEL=debug
Environment=KERNEL_LOG_MAX_SIZE=100M
Environment=KERNEL_LOG_ROTATE_COUNT=5

# Security settings
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
CapabilityBoundingSet=CAP_SYSLOG
AmbientCapabilities=CAP_SYSLOG

[Install]
WantedBy=multi-user.target